# roles/marathon/templates/marathon_from_src.conf.j2
# marathon.conf
# Generated by Ansible for {{ansible_fqdn}}

description "Marathon scheduler for Mesos"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

# Custom local
# Command line options: https://github.com/mesosphere/marathon#command-line-options
script
  exec nohup /opt/marathon/bin/start{% if groups.marathon|count > 1 %} --ha{% endif %} \
    --master {{zookeeper_url}}/mesos \
    --zk {{zookeeper_url}}/marathon \
    --hostname {{ mesos_local_address }} \
   >> /var/log/marathon.log 2>&1
end script
